# dist/cca-1.0.0-x86_64/install.sh
#!/usr/bin/env bash
#
# CCA Installation Script
# Installs CCA binaries and configuration
#

set -euo pipefail

# Colors
RED='\033[0;31m'
GREEN='\033[0;32m'
YELLOW='\033[1;33m'
BLUE='\033[0;34m'
NC='\033[0m'

# Defaults (can be overridden via environment or install.conf)
PREFIX="${PREFIX:-/usr/local}"
BIN_DIR="${BIN_DIR:-${PREFIX}/bin}"
CONFIG_DIR="${CONFIG_DIR:-${PREFIX}/etc/cca}"
DATA_DIR="${DATA_DIR:-${PREFIX}/share/cca}"
LOG_DIR="${LOG_DIR:-/var/log/cca}"

CCA_DAEMON_PORT="${CCA_DAEMON_PORT:-8580}"
CCA_ACP_PORT="${CCA_ACP_PORT:-8581}"
REDIS_PORT="${REDIS_PORT:-16379}"
POSTGRES_PORT="${POSTGRES_PORT:-15432}"
POSTGRES_USER="${POSTGRES_USER:-cca}"
POSTGRES_PASSWORD="${POSTGRES_PASSWORD:-cca_secure_password}"
POSTGRES_DB="${POSTGRES_DB:-cca}"

INSTALL_SYSTEMD="${INSTALL_SYSTEMD:-true}"
CREATE_USER="${CREATE_USER:-true}"
RUN_MIGRATIONS="${RUN_MIGRATIONS:-true}"

SCRIPT_DIR="$(cd "$(dirname "${BASH_SOURCE[0]}")" && pwd)"

# Load install.conf if exists
if [[ -f "${SCRIPT_DIR}/config/install.conf" ]]; then
    source "${SCRIPT_DIR}/config/install.conf"
elif [[ -f "${SCRIPT_DIR}/install.conf" ]]; then
    source "${SCRIPT_DIR}/install.conf"
fi

echo -e "${BLUE}╔════════════════════════════════════════════════════════════╗${NC}"
echo -e "${BLUE}║              CCA Installation Script                       ║${NC}"
echo -e "${BLUE}╚════════════════════════════════════════════════════════════╝${NC}"
echo ""
echo -e "  Install prefix:  ${GREEN}${PREFIX}${NC}"
echo -e "  Binaries:        ${GREEN}${BIN_DIR}${NC}"
echo -e "  Configuration:   ${GREEN}${CONFIG_DIR}${NC}"
echo -e "  Data directory:  ${GREEN}${DATA_DIR}${NC}"
echo -e "  Log directory:   ${GREEN}${LOG_DIR}${NC}"
echo ""

# Check for root
if [[ $EUID -ne 0 ]]; then
    echo -e "${RED}This script must be run as root (use sudo)${NC}"
    exit 1
fi

# Create CCA user if needed
if [[ "$CREATE_USER" == "true" ]]; then
    if ! id -u cca &>/dev/null; then
        echo -e "${YELLOW}Creating cca user...${NC}"
        useradd --system --no-create-home --shell /usr/sbin/nologin cca
    fi
fi

# Create directories
echo -e "${YELLOW}Creating directories...${NC}"
mkdir -p "$BIN_DIR" "$CONFIG_DIR" "$DATA_DIR" "$LOG_DIR" "$DATA_DIR/agents"

# Install binaries
echo -e "${YELLOW}Installing binaries...${NC}"
for bin in cca ccad cca-mcp; do
    if [[ -f "${SCRIPT_DIR}/bin/${bin}" ]]; then
        cp "${SCRIPT_DIR}/bin/${bin}" "${BIN_DIR}/"
        chmod +x "${BIN_DIR}/${bin}"
        echo -e "  ${GREEN}✓${NC} ${bin}"
    else
        echo -e "  ${RED}✗${NC} ${bin} (not found)"
    fi
done

# Generate API key
API_KEY=$(openssl rand -hex 32 2>/dev/null || head -c 32 /dev/urandom | xxd -p)

# Install configuration
echo -e "${YELLOW}Installing configuration...${NC}"
if [[ ! -f "${CONFIG_DIR}/cca.toml" ]]; then
    cat > "${CONFIG_DIR}/cca.toml" << TOMLEOF
# CCA Configuration
# Generated by installer

[daemon]
bind_address = "127.0.0.1:${CCA_DAEMON_PORT}"
log_level = "info"
max_agents = 10
log_file = "${LOG_DIR}/ccad.log"
data_dir = "${DATA_DIR}"
require_auth = true
api_keys = ["${API_KEY}"]

[redis]
url = "redis://localhost:${REDIS_PORT}"
pool_size = 10
context_ttl_seconds = 3600

[postgres]
url = "postgres://${POSTGRES_USER}:${POSTGRES_PASSWORD}@localhost:${POSTGRES_PORT}/${POSTGRES_DB}"
pool_size = 10
max_connections = 20

[agents]
default_timeout_seconds = 300
context_compression = true
token_budget_per_task = 50000
claude_path = "claude"

[acp]
websocket_port = ${CCA_ACP_PORT}
reconnect_interval_ms = 1000
max_reconnect_attempts = 5

[mcp]
enabled = true
bind_address = "127.0.0.1:9201"

[learning]
enabled = true
default_algorithm = "ppo"
training_batch_size = 32
update_interval_seconds = 300

[embeddings]
enabled = false
# ollama_url = "http://localhost:11434"
# model = "nomic-embed-text:latest"
# dimension = 768
TOMLEOF
    echo -e "  ${GREEN}✓${NC} cca.toml (generated)"
else
    echo -e "  ${YELLOW}⚠${NC} cca.toml (already exists, skipped)"
fi

# Create environment file
cat > "${CONFIG_DIR}/cca.env" << ENVEOF
# CCA Environment Variables
export CCA_CONFIG=${CONFIG_DIR}/cca.toml
export CCA_DATA_DIR=${DATA_DIR}
export CCA_DAEMON_URL=http://127.0.0.1:${CCA_DAEMON_PORT}
export CCA_ACP_URL=ws://127.0.0.1:${CCA_ACP_PORT}
export CCA_API_KEY=${API_KEY}
ENVEOF
echo -e "  ${GREEN}✓${NC} cca.env"

# Copy agent definitions
if [[ -d "${SCRIPT_DIR}/agents" ]]; then
    cp -r "${SCRIPT_DIR}/agents/"* "${DATA_DIR}/agents/" 2>/dev/null || true
    echo -e "  ${GREEN}✓${NC} Agent definitions"
fi

# Copy migrations
if [[ -d "${SCRIPT_DIR}/migrations" ]]; then
    mkdir -p "${DATA_DIR}/migrations"
    cp "${SCRIPT_DIR}/migrations/"*.sql "${DATA_DIR}/migrations/"
    echo -e "  ${GREEN}✓${NC} Database migrations"
fi

# Run migrations if PostgreSQL is available
if [[ "$RUN_MIGRATIONS" == "true" ]]; then
    echo -e "${YELLOW}Checking database...${NC}"
    if command -v psql &>/dev/null; then
        if PGPASSWORD="${POSTGRES_PASSWORD}" psql -h localhost -p "${POSTGRES_PORT}" -U "${POSTGRES_USER}" -d "${POSTGRES_DB}" -c "SELECT 1" &>/dev/null; then
            echo -e "  Running migrations..."
            for migration in "${DATA_DIR}/migrations/"*.sql; do
                PGPASSWORD="${POSTGRES_PASSWORD}" psql -h localhost -p "${POSTGRES_PORT}" -U "${POSTGRES_USER}" -d "${POSTGRES_DB}" -f "$migration" 2>&1 | grep -v "already exists" || true
            done
            echo -e "  ${GREEN}✓${NC} Migrations complete"
        else
            echo -e "  ${YELLOW}⚠${NC} PostgreSQL not reachable, skipping migrations"
        fi
    else
        echo -e "  ${YELLOW}⚠${NC} psql not found, skipping migrations"
    fi
fi

# Install systemd service
if [[ "$INSTALL_SYSTEMD" == "true" ]] && [[ -d /etc/systemd/system ]]; then
    echo -e "${YELLOW}Installing systemd service...${NC}"
    if [[ -f "${SCRIPT_DIR}/systemd/ccad.service" ]]; then
        cp "${SCRIPT_DIR}/systemd/ccad.service" /etc/systemd/system/
        systemctl daemon-reload
        echo -e "  ${GREEN}✓${NC} ccad.service installed"
        echo -e "  ${YELLOW}→${NC} Enable with: systemctl enable ccad"
        echo -e "  ${YELLOW}→${NC} Start with:  systemctl start ccad"
    fi
fi

# Set ownership
chown -R cca:cca "$DATA_DIR" "$LOG_DIR" 2>/dev/null || true
chmod 600 "${CONFIG_DIR}/cca.toml" "${CONFIG_DIR}/cca.env"

echo ""
echo -e "${GREEN}╔════════════════════════════════════════════════════════════╗${NC}"
echo -e "${GREEN}║              Installation Complete!                        ║${NC}"
echo -e "${GREEN}╚════════════════════════════════════════════════════════════╝${NC}"
echo ""
echo -e "  API Key: ${YELLOW}${API_KEY}${NC}"
echo ""
echo -e "  ${BLUE}Quick Start:${NC}"
echo -e "    1. Start Docker services (Redis + PostgreSQL):"
echo -e "       ${YELLOW}docker-compose up -d${NC}"
echo ""
echo -e "    2. Start the daemon:"
echo -e "       ${YELLOW}sudo systemctl start ccad${NC}"
echo -e "       or: ${YELLOW}ccad --config ${CONFIG_DIR}/cca.toml${NC}"
echo ""
echo -e "    3. Configure Claude Code MCP:"
echo -e "       Add to ~/.claude/mcp_servers.json:"
cat << 'MCPEOF'
       {
         "cca": {
           "command": "/usr/local/bin/cca-mcp",
           "args": [],
           "env": {
             "CCA_DAEMON_URL": "http://127.0.0.1:8580",
MCPEOF
echo -e "             \"CCA_API_KEY\": \"${API_KEY}\""
cat << 'MCPEOF'
           }
         }
       }
MCPEOF
echo ""
echo -e "  ${BLUE}Documentation:${NC} https://github.com/anthropics/cca"
echo ""
